# root/etc/systemd/system/pacman-auto-upgrade.service
[Unit]
Description=pacman-auto-upgrade

[Service]
ExecStart=/usr/bin/pacman -Syyu --noconfirm

[Install]
WantedBy=multi-user.target